# config/skel/.config/autostart/kill-fcitx.desktop
[Desktop Entry]
Exec=/usr/bin/kill-fcitx.sh
Name=kill-fcitx
OnlyShowIn=LXQt;
Type=Application
Version=1.0
X-LXQt-Need-Tray=true
